# daemon/claude-relay.service
[Unit]
Description=Ellie Relay — AI Processing + Multi-Channel I/O
After=network.target postgresql.service

[Service]
Type=simple
User=ellie
WorkingDirectory=/home/ellie/ellie-dev
ExecStart=/home/ellie/.bun/bin/bun run src/relay.ts
Restart=always
RestartSec=10

# Environment
Environment=HOME=/home/ellie
Environment=PATH=/home/ellie/.local/bin:/home/ellie/.bun/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/ellie/ellie-dev/.env

# Resource limits (ELLIE-240)
MemoryMax=3G
MemoryHigh=2G
CPUQuota=300%
TasksMax=512
OOMPolicy=kill

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ellie-relay

[Install]
WantedBy=multi-user.target
